# native/test/sv_service_standalone.service
[Unit]
Description=SV Subscriber Capture Service (Standalone C++ Backend)
Documentation=file:///home/powereureka/Desktop/subscriber/custom-sub/native/test/sv_service_standalone.cc
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/powereureka/Desktop/subscriber/custom-sub/native/test

# Binary path — update if needed
ExecStart=/home/powereureka/Desktop/subscriber/custom-sub/native/test/sv_service_standalone 0

# Restart on crash
Restart=on-failure
RestartSec=5s

# Output logging
StandardOutput=file:/var/log/sv_service/sv_service_standalone.log
StandardError=file:/var/log/sv_service/sv_service_error.log

# Resource limits
CPUAffinity=0-7
Nice=-5
MemoryLimit=2G
TasksMax=100

# Environment
Environment="LD_LIBRARY_PATH=/opt/intel/oneapi/mkl/latest/lib/intel64:/usr/lib"

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes

# Graceful shutdown
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=10s

[Install]
WantedBy=multi-user.target
